# service/cfs-repair.service
[Unit]
Description=ContainerFS Repair Service

[Service]
ExecStart=/usr/local/bin/cfs-repair /usr/local/conf/cfs-repair.ini
LimitNOFILE=1048576
LimitNPROC=1048576
LimitCORE=infinity
Restart=on-failure

[Install]
WantedBy=multi-user.target
